# DEBIAN/postinst
#!/bin/sh

systemctl start turboledz
systemctl enable turboledz

