# DEBIAN/prerm
#!/bin/sh

systemctl stop turboledz
systemctl disable turboledz

